# ansible/templates/etc/udev/rules.d/01-rtc.rules

ACTION=="add", KERNEL=="rtc0", RUN+="/sbin/hwclock --hctosys --utc", TAG+="systemd"
