# ansible/templates/etc/systemd/system/systemd-journald.service.d/wait_for_rtc.conf
[Unit]
Wants=dev-rtc0.device
After=dev-rtc0.device
